# deployment/systemd/mongoose-update.service
[Unit]
Description=mikro update service
OnFailure=status_email_user@%n.service

[Service]
ExecStart=%h/bin/update-mikro.sh
