# deployment/systemd/status_email_user@.service
[Unit]
Description=status email for %i to %u

[Service]
Type=oneshot
ExecStart=/usr/local/bin/systemd-email %u %i
Group=systemd-journal
